# dd_metal_backup/dd_backup.service
[Unit]
Description=Metal Backup Service

[Service]
Type=oneshot
ExecStart=/opt/dd_backup.sh

[Install]
WantedBy=multi-user.target
