# dd_metal_backup/dd_backup.timer
# Paste in /etc/systemd/system/dd_backup.timer
# Automatic running: sudo systemctl enable --now dd_backup.timer

[Unit]
Description=Metal Backup Service Timer

[Timer]
OnUnitActiveSec=3d
Persistent=true
Unit=dd_backup.service

[Install]
WantedBy=timers.target
