# dd_metal_backup/dd_backup.service
# Paste in /etc/systemd/system/dd_backup.service
# Automatic running: sudo systemctl enable --now dd_backup

[Unit]
Description=Metal Backup Service

[Service]
Type=oneshot
ExecStart=/opt/dd_backup.sh

[Install]
WantedBy=multi-user.target
